# supervisor/janus-sentinel.conf
[program:janus-sentinel]
command=/usr/local/bin/janus-sentinel
user=root
autostart=true
autorestart=true
stderr_logfile=/var/log/janus-sentinel.err.log
stdout_logfile=/var/log/janus-sentinel.out.log